# docker/restart-miner.cron
0 * * * * root /usr/bin/pkill socat
0 * * * * root /usr/bin/pkill banq
