# etc/systemd/system/banq-rtw.service
# /etc/systemd/system/banq-rt.service
[Unit]
Description=Run banq-rtw.sh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/banq-rtw.sh
