# etc/systemd/system/banq-riw.service
# /etc/systemd/system/banq-riw.service
[Unit]
Description=Run banq-riw.sh

[Service]
Type=oneshot
ExecStart=/usr/local/bin/banq-riw.sh
